# etc/local.d/update-mlocate-db.stop
#!/usr/bin/sh
echo "Updating mlocate database..."
/usr/bin/updatedb
